# config/includes.chroot/etc/systemd/system/betriebssystem-firstboot.service
[Unit]
Description=BETRIEBSSYSTEM first-boot setup (install Flathub apps)
After=network-online.target
Wants=network-online.target
# Run only on an INSTALLED system, never in the live session (which would try
# to download ~2 GB of Flatpaks on every boot from read-only media).
ConditionPathExists=!/run/live/medium
ConditionPathExists=!/var/lib/betriebssystem/firstboot.done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/mkdir -p /var/lib/betriebssystem
# Apps not in Debian and too heavy/namespace-bound to bake (Android Studio uses
# flatpak extra-data apply_extra). Logseq is NOT here -- it's baked into /opt by
# 0430-logseq so it works in the live session too.
ExecStart=/usr/bin/flatpak install -y --system --noninteractive flathub com.google.AndroidStudio cc.arduino.IDE2
ExecStartPost=/usr/bin/touch /var/lib/betriebssystem/firstboot.done

[Install]
WantedBy=multi-user.target
